# archarm/ondd/ondd.service
[Unit]
Description=Outernet Data Delivery agent service

[Service]
ExecStart=/usr/sbin/ondd

[Install]
WantedBy=multi-user.target
